# src/sysconf/udev/91-mydlp-usb.rules
SUBSYSTEMS=="usb", KERNEL=="sd??", ACTION=="add", RUN+="/usr/sbin/mydlp-mount-filterfs %k"
